# supervisor/zookeeper.conf
[program:zookeeper]
# command = /usr/bin/start-zookeeper.sh start-foreground
command = /tools/zookeeper/zookeeper-3.4.10/bin/zkServer.sh start-foreground
autostart = true
autorestart = true